# deploy/cryptobot.service
[Unit]
Description=Crypto Trading Bot (Donchian breakout on Alpaca)
Documentation=file:///home/trader/cryptobot/README.md
After=network-online.target mssql-server.service
Wants=network-online.target
Requires=mssql-server.service

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/home/trader/cryptobot
EnvironmentFile=/home/trader/cryptobot/.env
ExecStart=/home/trader/cryptobot/venv/bin/python -u main.py
Restart=always
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptobot

# Hardening (conservative — bot needs network + filesystem under WorkingDirectory)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/trader/cryptobot/logs /home/trader/cryptobot/backtest/results

[Install]
WantedBy=multi-user.target
